# deployment/vps/systemd/pharmacie-frontend.service.template
[Unit]
Description=Pharmacie Frontend Next.js
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=__APP_DIR__/frontend
EnvironmentFile=__APP_DIR__/frontend/.env.production
Environment=NODE_ENV=production
ExecStart=__NPM_BIN__ run start
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
