# deployment/vps/systemd/pharmacie-backend.service.template
[Unit]
Description=Pharmacie Backend API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=__APP_DIR__/backend
EnvironmentFile=__APP_DIR__/backend/.env.production
ExecStart=__NPM_BIN__ run start
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
